# installer/linux/systemd/vaybooks-bms.service
[Unit]
Description=VayBooks Boutique Management System
After=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=__INSTALL_DIR__/app
Environment=VAYBOOKS_DATA_DIR=__DATA_DIR__
ExecStart=__INSTALL_DIR__/python/python -m streamlit run app.py --server.headless=true --server.port=8501
Restart=on-failure
RestartSec=5
StandardOutput=append:__DATA_DIR__/logs/service.log
StandardError=append:__DATA_DIR__/logs/service.log

[Install]
WantedBy=multi-user.target
